# etc/profile.d/myenviornmentvariables.sh
export DJANGO_SETTINGS_MODULE=twizz.base_settings
